# 26_systemd/timer/chrisguest75capturephototimer.timer.template
[Unit]
Description=Capture a photo from webcam on a timer

[Timer]
OnCalendar=*-*-* *:*:00
Unit= chrisguest75capturephoto.service 
Persistent=true

[Install]
WantedBy=timers.target